# httpfrontend/conf/supervisor/httpserver.sv.conf
[program:httpserver]
command=/usr/local/sbin/httpserver
user = root
redirect_stderr=true
