# httpfrontend/conf/supervisor/nerve.sv.conf
[program:nerve]
command=/usr/local/bin/nerve -c /etc/nerve.conf.json
user = root
redirect_stderr=true
