# httpfrontend/conf/supervisor/nerve.sv.conf
[program:nerve]
command=/bin/bash -c "/usr/local/bin/expenv -i -f /etc/nerve.conf.json && /usr/local/bin/nerve -c /etc/nerve.conf.json"
user = root
redirect_stderr=true
