# httpfrontend/conf/supervisor/sshd.sv.conf
[program:sshd]
command=/usr/sbin/sshd -D
user = root
redirect_stderr=true