# httpfrontend/conf/supervisor/synapse.sv.conf
[program:synapse]
command=/bin/bash -c "/usr/local/bin/expenv -i -f /etc/synapse.conf.json && /usr/local/bin/synapse -c /etc/synapse.conf.json"
user = root
redirect_stderr=true