# service/conf/supervisor/httpservice.sv.conf
[program:httpservice]
command=/usr/local/sbin/httpservice
user = root
redirect_stderr=true
