# service/conf/supervisor/httpservice.sv.conf
[program:httpservice]
command=/bin/bash -c "/usr/local/sbin/httpservice -name=$SMARTS_INSTANCE_ID"
user = root
redirect_stderr=true
