# service/conf/supervisor/synapse.sv.conf
[program:synapse]
command=/usr/local/bin/synapse -c /etc/synapse.conf.json
user = root
redirect_stderr=true